# launchd/com.mcurtis.project-tracker-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Prompt bridge launchd agent (Phase 14).

  Auto-starts scripts/prompt-bridge.js at login and restarts it on crash.
  The two __REPO__ placeholders + the node path are rewritten by the
  install snippet in the README before bootstrapping into launchd.

  The bridge binds 127.0.0.1:4318 only — loopback, never network — and
  spawns `claude` in the requested cwd when the browser prompts it.

  Install / uninstall: see README, section "Daemon".
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mcurtis.project-tracker-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>__REPO__/scripts/prompt-bridge.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:__HOME__/.local/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>__REPO__/data/bridge.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO__/data/bridge.err</string>
</dict>
</plist>
